# examples/demo/config/supervisor.conf
[supervisord]
nodaemon=true

[program:app]
command = python3 /app/app.py

[program:sshd]
command=/usr/sbin/sshd -D
